# amarelo-earn.service
[Unit]
Description=Amarelo Earn - Mega-Sena Static Page Generator
After=network.target

[Service]
Type=oneshot
User=roberto
WorkingDirectory=/home/roberto/Projetos/earn
ExecStart=/tmp/amarelo_venv/bin/python /home/roberto/Projetos/earn/generate.py

[Install]
WantedBy=multi-user.target
